# adjwin.timer
[Unit]
Description = Run the remindme timer once a minute

[Timer]
OnCalendar = *:*:0/2
AccuracySec = 1s

[Install]
WantedBy = timers.target
